# build/start-scripts/galaxy-agent
#!/bin/bash -l
# chkconfig: 2345 80 90
# description: Activates/Deactivates Galaxy agent

GALAXY_CONFIG="/etc/galaxy.conf"
GALAXY_USER=`ruby -ryaml -e "puts YAML.load(File.open('$GALAXY_CONFIG'))['galaxy.agent.user'] || ''"`

if [ -z "$GALAXY_USER" ]
then
    echo "Error: Unable to determine galaxy agent user. Please set the galaxy.agent.user property in $GALAXY_CONFIG."
    exit 2
fi

case "$1" in
'start')
        logger -p daemon.notice "Starting galaxy-agent"
        su - $GALAXY_USER -c 'galaxy-agent --start'
        ;;
'stop')
        logger -p daemon.notice "Stopping galaxy-agent"
        su - $GALAXY_USER -c 'galaxy-agent --stop'
        ;;
'restart')
        logger -p daemon.notice "Restarting galaxy-agent"
        su - $GALAXY_USER -c 'galaxy-agent --stop'
        su - $GALAXY_USER -c 'galaxy-agent --start'
        ;;
'status') 
        galaxy-agent --status
        ;;
*)
        echo "Usage: $0 <start|stop|restart|status>"
        exit 1
        ;;
esac

